# etc/systemd/system/disk-monitor.service
[Unit]
Description=Disk Usage Monitor

[Service]
Type=oneshot
ExecStart=/usr/local/bin/disk_monitor.sh
