# etc/systemd/system/log-watcher.service
[Unit]
Description=Log Watcher Service

[Service]
ExecStart=/usr/local/bin/log_watcher.sh
Restart=always
RestartSec=60
