# etc/systemd/system/ram-monitor.service
[Unit]
Description=RAM Usage Monitor

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ram_monitor.sh
